# scripts/scanoss-dependencies-api.service
[Unit]
Description=SCANOSS Dependencies API Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=scanoss
ExecStart=/usr/local/bin/scanoss-dependencies-api.sh

[Install]
WantedBy=multi-user.target
